# deploy/systemd/polymonitor-backend.service
[Unit]
Description=PolyMonitor FastAPI backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=polymonitor
Group=polymonitor
WorkingDirectory=/opt/polymonitor/backend
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/polymonitor/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
