# launchd/com.polymonitor.backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.polymonitor.backend</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/infinity/Python/Poly_Codex/scripts/start_backend.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/infinity/Python/Poly_Codex/backend</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/infinity/Python/Poly_Codex/logs/backend.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/infinity/Python/Poly_Codex/logs/backend.err.log</string>
</dict>
</plist>
