# launchd/com.polymonitor.remote.backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.polymonitor.remote.backend</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/roc/Poly_Codex/scripts/remote_start_backend.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/roc/Poly_Codex/backend</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/roc/Poly_Codex/logs/remote-backend.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/roc/Poly_Codex/logs/remote-backend.err.log</string>
</dict>
</plist>
